# debian/prerm
#!/bin/sh

set -e

if [ -x "/etc/init.d/jabclnd" ]; then
	invoke-rc.d jabclnd stop || exit $?
fi

#DEBHELPER#
